# debian/ndhcpd.service
[Unit]
Description=Nano-DHCP server
After=network.target

[Service]
ExecStart=/usr/bin/ndhcpd
Type=forking

[Install]
WantedBy=default.target
